# systemd/mysqlchk.socket
[Unit]
Description=Percona XtraDB Cluster node check socket

[Socket]
ListenStream=9200
Accept=yes

[Install]
WantedBy=sockets.target