# systemd/mysqlchk@.service
[Unit]
Description=Percona XtraDB Cluster node check service

[Service]
ExecStart=-/usr/bin/clustercheck
StandardInput=socket